# examples/cron.d/__console

# Файл из репозитория MikBill-DaemonSystem-PHP-Kit
# https://github.com/itpanda-llc/mikbill-daemonsystem-php-kit

#do_api_terminal_payments
* * * * * root cd /var/www/mikbill/admin/ && /usr/bin/php ./index.php "do_api_terminal_payments" > /dev/null 2>&1
